# com.gemini-proxy.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gemini-proxy</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>__BINARY_PATH__</string>
        <string>-use-keychain</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>__WORKING_DIR__</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/gemini-proxy.log</string>
    
    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/gemini-proxy.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PORT</key>
        <string>9877</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
</dict>
</plist>